# sqm_reader.service
[Service]
Type=oneshot
ExecStart=/usr/bin/echo 'serial SQM reader'
StopWhenUnneeded=true
